# bootc/common/tailscale/usr/lib/systemd/system/tailscale-register.service
[Unit]
Description=Tailscale Registration (Doppler)
After=tailscaled.service network-online.target
Requires=tailscaled.service network-online.target
ConditionPathExists=/etc/doppler.env

[Service]
Type=oneshot
Environment=HOME=/root
EnvironmentFile=/etc/doppler.env
# Only run registration if tailscale is not already logged in/online
ExecStart=/usr/bin/sh -c 'tailscale status >/dev/null 2>&1 || /usr/bin/doppler run -- tailscale up --authkey=$$TS_AUTHKEY --hostname=%H'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
